# wireguard-webconfig.service
[Unit]
Description=Wireguard Web Configuration
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/usr/bin/wireguard-webconfig /etc/wireguard/wg0.conf

[Install]
WantedBy=multi-user.target
